# s3-mount.service
[Unit]
Description=Mount S3 Bucket for Virtual City Tours
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="S3_BUCKET=your-bucket-name"
WorkingDirectory=/path/to/virtueller-reisefuehrer
ExecStart=/path/to/virtueller-reisefuehrer/setup-s3.sh
ExecStop=/path/to/virtueller-reisefuehrer/unmount-s3.sh

[Install]
WantedBy=multi-user.target
